# crontab
* * * 1 1     /root/proxmox-stuff/prox_config_backup.sh
