# etc/vfxws-013/rc.local
#!/bin/sh

# pull latest setup from github
CD=$(pwd)
if [ -e /root/proxmox-stuff ] ; then
	cd /root/proxmox-stuff
	git reset HEAD~5 --hard
	while true ; do
		git pull
	 	if [ $? -ne 0 ] ; then
			break
		fi
		echo "Failed to pull from github... retrying..."
		sleep 1
	done
	./setup.py -i
	echo $?
	cd $CD
fi


pvecm expected 1
#/usr/share/perl5/PVE/QemuServer.pm

#/root/testStart

# start windows
# startvm 401

# start flame
startFlame(){
	sleep 180
	startvm 700
}
startFlame &

exit 0
