# etc/vfxws-013/rc.local
#!/bin/sh


pvecm expected 1
#/usr/share/perl5/PVE/QemuServer.pm

#/root/testStart
#startvm 401
startvm 700

exit 0
