# etc/centos/rc.local
#!/bin/bash



blockdev --setra 131072 /dev/sd?




exit 0
