# etc/manjaro/rc.local
#!/bin/bash

# pull latest setup from github
CD=$(pwd)
if [ -e /root/proxmox-stuff ] ; then
	cd /root/proxmox-stuff
	git reset HEAD~5 --hard
	while true ; do
		git pull
	 	if [ $? == 0 ] ; then
			break
		fi
		echo "Failed to pull from github... retrying..."
		sleep 1
	done
	./setup.py -i
	echo $?
	cd $CD
fi



exit 0
